# systemd/gtfs-update.timer
[Unit]
Description=Run GTFS Auto-Update every 15 days

[Timer]
# Run every 15 days at 4:00 AM
OnCalendar=*-*-01,16 04:00:00
# Also run 15 minutes after boot if last run was missed
Persistent=true

[Install]
WantedBy=timers.target
