# systemd/gtfs-update.timer
[Unit]
Description=Run GTFS Auto-Update weekly

[Timer]
# Run every Monday at 4:00 AM
OnCalendar=Mon 04:00:00
# Also run 15 minutes after boot if last run was missed
Persistent=true

[Install]
WantedBy=timers.target
